# .docker/supervisord.conf
[program:nginx]
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf
[program:php-fpm]
command=/usr/local/sbin/php-fpm
[supervisord]
logfile=/var/log/supervisord.log 
nodaemon=true